# supervisord.conf
[supervisord]
nodaemon=true

[program:postgres]
user=postgres
command=/usr/lib/postgresql/9.5/bin/postgres -D /var/lib/postgresql/9.5/main -c config_file=/etc/postgresql/9.5/main/postgresql.conf

[program:php5-fpm]
command=/usr/sbin/php5-fpm

[program:nginx]
command=/usr/sbin/nginx
